# cfg/roles/home-backups/files/borg-backup.target
[Unit]
Description=Timer target for Borg backup command
StopWhenUnneeded=yes
